# cp_sim201.service
[Unit]
Description=CP700P EV Charger OCPP 2.0.1 Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/home/jeongsooh/cp_sim201/scripts/start_daemon.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
